# deploy/starling/aeroprint-chatter.service
[Unit]
Description=px4
After=docker-autorun.service
Requires=docker-autorun.service

[Service]
StartLimitIntervalSec=30
StartLimitBurst=5


WorkingDirectory=/home/root/aeroprint/deploy/starling
ExecStart=/usr/bin docker compose up
Restart=on-failure
RestartSec=5s


[Install]
WantedBy=multi-user.target
